# files/system/etc/systemd/system/ammix-bootc-update.service
[Unit]
Description=Ammix Custom: Bootc Upgrade
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc upgrade
